# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true

[program:vsftpd]
command=/usr/sbin/vsftpd
autostart=true
autorestart=true

[program:samba]
command=smbd -s -S -F /etc/samba/smb.conf
autostart=true
autorestart=true